# build/wg-watchdog.service
[Unit]
Description=WireGuard DDNS watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wg-watchdog
